# systemd/conflux@.system.service
# System-wide conflux daemon, profile-templated.
#
# The systemd instance name selects the profile: `conflux@desktop` runs only
# the syncs whose `profiles` include "desktop". This lets one shared config
# file drive different profiles on different hosts. To run the implicit
# "default" profile, use the plain (non-templated) conflux.system.service.
#
# Install:
#   sudo cp target/release/conflux /usr/local/bin/conflux
#   sudo useradd --system --create-home --home-dir /var/lib/conflux conflux
#   sudo cp systemd/conflux@.system.service /etc/systemd/system/conflux@.service
#   sudo cp config.example.toml /etc/conflux/config.toml   # then edit
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now conflux@desktop
#
# Multiple instances can run at once: each profile gets its own state
# subdirectory (/var/lib/conflux/<profile>) and control socket
# (/run/conflux/conflux-<profile>.sock). Target one with, e.g.,
# `conflux --system --profile desktop status`.

[Unit]
Description=conflux file-sync daemon (system, profile %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=CONFLUX_LOG=info
ExecStart=/usr/local/bin/conflux --system --profile %i daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
User=conflux
Group=conflux
RuntimeDirectory=conflux
StateDirectory=conflux
ConfigurationDirectory=conflux

[Install]
WantedBy=multi-user.target
